# pkginfo/Scripts/postinstall
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin; export PATH

# ------------------------------------------------------------
# relaunch karabiner_grabber

killall karabiner_grabber
killall karabiner_console_user_server

launchctl bootstrap system /Library/LaunchDaemons/org.pqrs.karabiner.karabiner_grabber.plist
launchctl enable system/org.pqrs.karabiner.karabiner_grabber

exit 0
