# pkginfo/Scripts/postinstall
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin; export PATH

# ------------------------------------------------------------
# relaunch karabiner_grabber

killall karabiner_grabber

launchctl load -w /Library/LaunchDaemons/org.pqrs.Karabiner-Elements.plist

exit 0
